# supervisord.conf
[supervisord]
nodaemon=True


[program:controller]
command=/usr/local/bin/python -m memes_service

[program:worker]
command=/usr/local/bin/python -m memes_service/worker
